# apps/agents/Worker/france-nuage-agent/DEBIAN/postinst
#!/bin/bash
set -e

# Create a new group and user for the service if they don't exist
if ! getent group france-nuage-agent > /dev/null 2>&1; then
    groupadd -r france-nuage-agent
fi

if ! getent passwd france-nuage-agent > /dev/null 2>&1; then
    useradd -r -s /usr/sbin/nologin -g france-nuage-agent france-nuage-agent
fi
echo "User and group created"

# Create the log directory and set permissions
LOG_DIR="/var/log/france-nuage-agent"
LOG_FILE="$LOG_DIR/metrics.logs"
mkdir -p "$LOG_DIR"
touch "$LOG_FILE"
chown france-nuage-agent:france-nuage-agent "$LOG_DIR" "$LOG_FILE"
chmod 755 "$LOG_DIR"
chmod 644 "$LOG_FILE"

cd /opt/france-nuage-agent

# Ensure correct permissions
chown -R france-nuage-agent:france-nuage-agent /opt/france-nuage-agent

# Remove existing node_modules and package-lock.json
rm -rf node_modules package-lock.json

# Clear npm cache
npm cache clean --force

# Reinstall dependencies
if ! npm install --legacy-peer-deps; then
    echo "npm install failed"
    exit 1
fi

# Compile TypeScript files
if ! npx tsc; then
    echo "TypeScript compilation failed"
    exit 1
fi

# Reload systemd manager configuration
systemctl daemon-reload

# Enable and start the service
systemctl enable france-nuage-agent.service
systemctl start france-nuage-agent.service
echo "Post-installation script completed successfully"
